# gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/coursemanager-gunicorn.sock

[Install]
WantedBy=sockets.target